# supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; Kronos Trading System — Supervisord Configuration
; Target: Linux Ubuntu VPS (Hostinger Singapore, 4 GB RAM)
;
; Deploy checklist:
;   1. Copy entire project to /app/kronos/
;   2. Create virtualenv: python3 -m venv /app/venv
;   3. Install deps:      /app/venv/bin/pip install -r requirements.txt
;   4. Copy .env.example → /app/kronos/.env, fill in all secrets
;   5. Create log dir:    sudo mkdir -p /var/log/kronos && sudo chown kronos:kronos /var/log/kronos
;   6. Create data dir:   mkdir -p /app/kronos/data/reports
;   7. Create models dir: mkdir -p /app/kronos/models
;   8. Run once to init DB: cd /app/kronos && /app/venv/bin/python -c "from db import init_db; init_db()"
;   9. Start:             sudo supervisord -c /app/kronos/supervisord.conf
;  10. Status:            sudo supervisorctl -c /app/kronos/supervisord.conf status
;
; SECURITY: API keys are loaded from /app/kronos/.env via the environment directive.
;           Never commit .env to version control.
;
; Process names match Module 11 auto-restart expectations:
;   kronos-data-collection   (Module 1)
;   kronos-slippage-model    (Module 2)
;   kronos-signal-generator  (Module 4)
;   kronos-risk-check        (Module 5)
;   kronos-execution         (Module 6)
;   kronos-position-monitor  (Module 7)
;   kronos-portfolio-manager (Module 8)
;   kronos-tax-tracker       (Module 9)
;   kronos-notifier          (Module 10)
;   kronos-health-monitor    (Module 11)
;
; Module 3 (macro_calendar.py) is a library — no standalone process.
; ─────────────────────────────────────────────────────────────────────────────

[unix_http_server]
file=/tmp/kronos-supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/kronos/supervisord.log
logfile_maxbytes=20MB
logfile_backups=5
loglevel=info
nodaemon=false
pidfile=/tmp/kronos-supervisord.pid
user=kronos

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/kronos-supervisor.sock

; ─────────────────────────────────────────────────────────────────────────────
; Shared environment — loaded by all programs
; Secrets are read from /app/kronos/.env at supervisord start time.
; Add any new env vars here so all modules see them.
; ─────────────────────────────────────────────────────────────────────────────

; NOTE: supervisord does not natively source .env files. Use one of:
;   Option A (recommended): set secrets via system /etc/environment or
;                            a systemd EnvironmentFile before starting supervisord.
;   Option B: inline secrets below (NOT recommended — visible in process list).
;
; The environment below lists all non-secret vars with their defaults.
; Override any of these in your systemd unit or login profile before starting.

[program:kronos-data-collection]
command=/app/venv/bin/python /app/kronos/01_data_collection.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/01_data_collection.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/01_data_collection.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_API_KEY="%(ENV_KRONOS_API_KEY)s",
    KRONOS_API_SECRET="%(ENV_KRONOS_API_SECRET)s"

[program:kronos-slippage-model]
command=/app/venv/bin/python /app/kronos/02_slippage_model.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/02_slippage_model.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/02_slippage_model.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s"

[program:kronos-signal-generator]
command=/app/venv/bin/python /app/kronos/04_signal_generator.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/04_signal_generator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/04_signal_generator.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_MODEL_PATH="/app/kronos/models/kronos_model.pt",
    KRONOS_SEQ_LEN="96",
    KRONOS_PRED_LEN="6",
    KRONOS_MAX_PRED_RETURN="20.0",
    KRONOS_MACRO_CALENDAR="/app/kronos/macro_calendar.json",
    HF_HOME="/app/kronos/.cache/huggingface",
    HF_HUB_DISABLE_IMPLICIT_TOKEN="1"
command=/app/venv/bin/python /app/kronos/05_risk_check.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/05_risk_check.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/05_risk_check.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_MACRO_CALENDAR="/app/kronos/macro_calendar.json",
    KRONOS_CB_SPREAD_PCT="1.0"

[program:kronos-execution]
command=/app/venv/bin/python /app/kronos/06_execution.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=60
stdout_logfile=/var/log/kronos/06_execution.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/06_execution.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_API_KEY="%(ENV_KRONOS_API_KEY)s",
    KRONOS_API_SECRET="%(ENV_KRONOS_API_SECRET)s",
    KRONOS_LEVERAGE="2.0",
    KRONOS_USD_INR_RATE="%(ENV_KRONOS_USD_INR_RATE)s",
    KRONOS_PORTFOLIO_VALUE_INR="100000.0"

[program:kronos-position-monitor]
command=/app/venv/bin/python /app/kronos/07_position_monitor.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=60
stdout_logfile=/var/log/kronos/07_position_monitor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/07_position_monitor.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_API_KEY="%(ENV_KRONOS_API_KEY)s",
    KRONOS_API_SECRET="%(ENV_KRONOS_API_SECRET)s",
    KRONOS_USD_INR_RATE="%(ENV_KRONOS_USD_INR_RATE)s"

[program:kronos-portfolio-manager]
command=/app/venv/bin/python /app/kronos/08_portfolio_manager.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
; Longer stop wait — in-flight position close orders should complete
stopwaitsecs=120
stdout_logfile=/var/log/kronos/08_portfolio_manager.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/08_portfolio_manager.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_API_KEY="%(ENV_KRONOS_API_KEY)s",
    KRONOS_API_SECRET="%(ENV_KRONOS_API_SECRET)s",
    KRONOS_STARTING_CAPITAL_INR="100000.0",
    KRONOS_USD_INR_RATE="%(ENV_KRONOS_USD_INR_RATE)s",
    KRONOS_PHASE="%(ENV_KRONOS_PHASE)s",
    KRONOS_MONTHLY_FIXED_COSTS_INR="915.0"

[program:kronos-tax-tracker]
command=/app/venv/bin/python /app/kronos/09_tax_tracker.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/09_tax_tracker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/09_tax_tracker.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_USD_INR_RATE="%(ENV_KRONOS_USD_INR_RATE)s",
    KRONOS_STARTING_CAPITAL_INR="100000.0"

[program:kronos-notifier]
command=/app/venv/bin/python /app/kronos/10_notification.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=15
; Telegram bot needs longer startup time (downloads bot info)
stopwaitsecs=30
stdout_logfile=/var/log/kronos/10_notification.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/10_notification.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_TELEGRAM_BOT_TOKEN="%(ENV_KRONOS_TELEGRAM_BOT_TOKEN)s",
    KRONOS_TELEGRAM_CHAT_ID="%(ENV_KRONOS_TELEGRAM_CHAT_ID)s",
    KRONOS_EMAIL_FROM="%(ENV_KRONOS_EMAIL_FROM)s",
    KRONOS_EMAIL_TO="%(ENV_KRONOS_EMAIL_TO)s",
    KRONOS_EMAIL_APP_PASSWORD="%(ENV_KRONOS_EMAIL_APP_PASSWORD)s",
    KRONOS_NOTIFIER_STATE_PATH="/app/kronos/data/notifier_state.json",
    KRONOS_STARTING_CAPITAL_INR="100000.0"

[program:kronos-health-monitor]
command=/app/venv/bin/python /app/kronos/11_health_monitor.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/11_health_monitor.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/11_health_monitor.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_NOTIFIER_STATE_PATH="/app/kronos/data/notifier_state.json",
    KRONOS_SUPERVISORD_PREFIX="kronos-"

; ─────────────────────────────────────────────────────────────────────────────
; Module 12 — Dashboard (Flask, served on KRONOS_DASHBOARD_PORT, default 8050)
; ─────────────────────────────────────────────────────────────────────────────

[program:kronos-dashboard]
command=/app/venv/bin/python /app/kronos/dashboard.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=30
stdout_logfile=/var/log/kronos/12_dashboard.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/12_dashboard.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_DASHBOARD_PORT="%(ENV_KRONOS_DASHBOARD_PORT)s",
    KRONOS_DASHBOARD_USER="%(ENV_KRONOS_DASHBOARD_USER)s",
    KRONOS_DASHBOARD_PASS="%(ENV_KRONOS_DASHBOARD_PASS)s",
    KRONOS_PHASE="%(ENV_KRONOS_PHASE)s",
    KRONOS_STARTING_CAPITAL_INR="100000.0"

; ─────────────────────────────────────────────────────────────────────────────
; Module 15 — Kronos-mini-4h signal generator (4H, 24H horizon)
; ─────────────────────────────────────────────────────────────────────────────

[program:kronos-mini-4h-generator]
command=/app/venv/bin/python /app/kronos/15_mini_4h_generator.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=30
stopwaitsecs=60
stdout_logfile=/var/log/kronos/15_mini_4h_generator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/15_mini_4h_generator.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_SHADOW_DEVICE="cpu",
    CUDA_VISIBLE_DEVICES="-1",
    KRONOS_SHADOW_SAMPLE_COUNT="%(ENV_KRONOS_SHADOW_SAMPLE_COUNT)s",
    KRONOS_MINI_4H_CONTEXT="2048",
    HF_HUB_DISABLE_IMPLICIT_TOKEN="1",
    HF_HOME="/app/kronos/.cache/huggingface",
    PYTHONUNBUFFERED="1"

; ─────────────────────────────────────────────────────────────────────────────
; Module 16 — Kronos-base-4h signal generator (4H, 24H horizon)
; ─────────────────────────────────────────────────────────────────────────────

[program:kronos-base-4h-generator]
command=/app/venv/bin/python /app/kronos/16_base_4h_generator.py
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=30
stopwaitsecs=60
stdout_logfile=/var/log/kronos/16_base_4h_generator.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
stderr_logfile=/var/log/kronos/16_base_4h_generator.err
stderr_logfile_maxbytes=10MB
environment=
    KRONOS_DB_PATH="/app/kronos/data/kronos.db",
    KRONOS_PAPER_MODE="%(ENV_KRONOS_PAPER_MODE)s",
    KRONOS_SHADOW_DEVICE="cpu",
    CUDA_VISIBLE_DEVICES="-1",
    KRONOS_SHADOW_SAMPLE_COUNT="%(ENV_KRONOS_SHADOW_SAMPLE_COUNT)s",
    KRONOS_BASE_4H_CONTEXT="512",
    HF_HUB_DISABLE_IMPLICIT_TOKEN="1",
    HF_HOME="/app/kronos/.cache/huggingface",
    PYTHONUNBUFFERED="1"

; ─────────────────────────────────────────────────────────────────────────────
; Process group — allows: supervisorctl restart kronos:*
; ─────────────────────────────────────────────────────────────────────────────

[group:kronos]
programs=kronos-data-collection,kronos-slippage-model,kronos-signal-generator,kronos-risk-check,kronos-execution,kronos-position-monitor,kronos-portfolio-manager,kronos-tax-tracker,kronos-notifier,kronos-health-monitor,kronos-dashboard,kronos-mini-4h-generator,kronos-base-4h-generator,kronos-cloudflared

; ── Cloudflare Tunnel — expose dashboard :8050 ───────────────────────────────
; Uses a quick tunnel (no account needed). URL changes on each restart.
; Find the URL: grep -i "trycloudflare.com" /var/log/kronos/cloudflared.log
[program:kronos-cloudflared]
command=/usr/local/bin/cloudflared tunnel --url http://localhost:8050
directory=/app/kronos
user=kronos
autostart=true
autorestart=true
startretries=5
startsecs=5
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/kronos/cloudflared.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
